# .config/systemd/user/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15m
Unit=offlineimap.service

[Install]
WantedBy=mail.target
